# maestro/rootfs/etc/services.d/maestro/run
#!/usr/bin/with-contenv bashio
# Start Hass Maestro Swift server
exec /usr/bin/maestro
